# roles/superset/templates/etc/profile.d/superset.sh
export PYTHONPATH="${PYTHONPATH}:{{ superset_home }}"